# crontab
# 抓取脚本
*/30 * * * * cd /var/www/pythonzone/web; python ./manage.py crawl >> crawl.log 2>&1
*/1 * * * * cd /var/www/pythonzone/web; python ./manage.py crawl_topics >> crawl_topics.log 2>&1